# door_access.service
[Unit]
Description=Door Access Control System
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=door
Group=door
ExecStart=/usr/bin/python3 /opt/door_access/main.py
# rpi-lgpio/lgpio creates a notify FIFO (.lgd-nfy*) in the CWD, so the
# service needs a writable working directory. RuntimeDirectory creates
# /run/door_access (owned by the service user) before ExecStart.
RuntimeDirectory=door_access
RuntimeDirectoryMode=0750
WorkingDirectory=/run/door_access
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=door_access
Environment=PYTHONUNBUFFERED=1
KillMode=mixed
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
